# config/reqpmin.conf
description "Make sure reqmin is running"

start on local-filesystem
stop on shutdown

respawn 
exec /etc/check_rpm start
